# node-red.conf
#!upstart
description "node-red"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

exec node /var/www/node-red/red.js >> /var/log/node-red.log 2>&1
